# deploy/incudal-online-update@.service.example
[Unit]
Description=PayIncus online update task %i
Documentation=https://github.com/VipMaxxxx/payincus
After=network.target postgresql.service redis-server.service
Requires=postgresql.service redis-server.service

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/incudal
EnvironmentFile=/opt/incudal/.env
Environment=HOME=/var/cache/incudal-ota
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/local/libexec/incudal/incudal-online-task update %i
TimeoutStartSec=1800
StandardOutput=journal
StandardError=journal
SyslogIdentifier=incudal-online-update
